# supervisor.conf
[program:sesha3]
command=/home/ubuntu/sesha3/sesha3 --syslog
directory=/home/ubuntu/sesha3
autostart=true
autorestart=true
stderr_logfile=syslog
stdout_logfile=syslog
