# supervisor.conf
[program:sesha3]
command=/home/ec2-user/sesha3/sesha3 serve --syslog
directory=/home/ec2-user/sesha3
autostart=true
autorestart=true
stderr_logfile=syslog
stdout_logfile=syslog
